# deploy/launchd/com.local.invest-agent-advisor-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.invest-agent-advisor-scheduler</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/apple/Documents/AIinvestmentAgent/.venv/bin/python</string>
    <string>-m</string>
    <string>invest_agent.cli</string>
    <string>advisor-scheduler-loop</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/apple/Documents/AIinvestmentAgent</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/apple/Documents/AIinvestmentAgent/logs/invest-agent-advisor-scheduler.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/apple/Documents/AIinvestmentAgent/logs/invest-agent-advisor-scheduler.err.log</string>
</dict>
</plist>
