# deploy/launchd/com.local.invest-agent-advisor-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.local.invest-agent-advisor-scheduler</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/apple/Documents/AIinvestmentAgent &amp;&amp; exec /Users/apple/Documents/AIinvestmentAgent/.venv/bin/python -m invest_agent.cli advisor-scheduler-loop</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/apple/Documents/AIinvestmentAgent</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/apple/Documents/AIinvestmentAgent/.venv/bin:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>SoftResourceLimits</key>
  <dict>
    <key>NumberOfFiles</key>
    <integer>4096</integer>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/apple/Documents/AIinvestmentAgent/logs/invest-agent-advisor-scheduler.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/apple/Documents/AIinvestmentAgent/logs/invest-agent-advisor-scheduler.err.log</string>
</dict>
</plist>
